# sdcard/rootfs/etc/systemd/system/pwnagotchi.service
[Unit]
Description=Pwnagotchi Autostart

[Service]
Type=forking
ExecStart=/root/pwnagotchi/scripts/startup.sh
Restart=always

[Install]
WantedBy=multi-user.target
